# qr3_49.conf
[program:qr3_49]
command=python3 /home/curator3/qr3_49/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator3/qr3_49/test.err.log
stdout_logfile=/home/curator3/qr3_49/test.err.log
